# test.service
[Unit]
Description=test

[Service]
User=root
WorkingDirectory=/root
ExecStart=/bin/bash -c 'nvidia-smi -L'

[Install]
WantedBy=multi-user.target
